# etc/profile.d/armitage.sh.j2
if [ -d /opt/armitage ] ; then
  PATH=$PATH:/opt/armitage
fi
