# etc/profile.d/motd.sh.j2
if [ -r /etc/motd.logo ] ; then
  cat /etc/motd.logo
fi
