# etc/profile.d/msf.sh.j2
if [ -d /opt/metasploit-framework ] ; then
  PATH=$PATH:/opt/metasploit-framework
  export MSF_DATABASE_CONFIG=/opt/metasploit-framework/database.yml
fi
