# etc/profile.d/rvm-ps1.sh.j2
if [ -r "/usr/local/rvm/contrib/ps1_functions" ] ; then
  source "/usr/local/rvm/contrib/ps1_functions"
  ps1_set --prompt '[\\h]\ \\W\ \(\$?\)\\\$'
fi
