# etc/init/start-mpstat.conf.j2
start on stopped rc RUNLEVEL=[2345]

env MPSTAT_TTY={{ttys.mpstat}}
task
script
. /etc/sysconfig/init
initctl start mpstat TTY=$MPSTAT_TTY
end script
