# etc/init/start-gstat.conf.j2
start on stopped rc RUNLEVEL=[2345]

env GSTAT_TTY={{ttys.gstat}}
task
script
. /etc/sysconfig/init
initctl start gstat TTY=$GSTAT_TTY
end script
