# etc/init/start-messages.conf.j2
start on stopped rc RUNLEVEL=[2345]

env MESSAGES_TTY={{ttys.messages}}
task
script
. /etc/sysconfig/init
initctl start messages TTY=$MESSAGES_TTY
end script
