# etc/init/start-iostat.conf.j2
start on stopped rc RUNLEVEL=[2345]

env IOSTAT_TTY={{ttys.iostat}}
task
script
. /etc/sysconfig/init
initctl start iostat TTY=$IOSTAT_TTY
end script
