# etc/init/start-top-mem.conf.j2
start on stopped rc RUNLEVEL=[2345]

env TOP_MEM_TTY={{ttys.top_mem}}
task
script
. /etc/sysconfig/init
initctl start top-mem TTY=$TOP_MEM_TTY
end script
